# deploy/ai-resume.service
[Unit]
Description=AI Resume Generator (Gunicorn)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/ai-resume
Environment="PATH=/opt/ai-resume/venv/bin"
EnvironmentFile=/opt/ai-resume/.env
ExecStart=/opt/ai-resume/venv/bin/gunicorn wsgi:app \
    --workers 1 \
    --threads 2 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --max-requests 500 \
    --max-requests-jitter 50 \
    --worker-tmp-dir /dev/shm \
    --access-logfile /var/log/ai-resume/access.log \
    --error-logfile /var/log/ai-resume/error.log \
    --log-level warning

# Restart on failure
Restart=always
RestartSec=5

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
